# api/healtcheck.cron
*/10 * * * * bash -c '/opt/healtcheck.sh'